# templates/systemd/jobbored-refresh.service
[Unit]
Description=JobBored daily discovery refresh
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory={{REPO_ROOT}}
ExecStart={{NODE_PATH}} {{SCRIPT_PATH}} --trigger scheduled-local --port {{PORT}} --sheet-id {{SHEET_ID}}
