# templates/launchd/com.jobbored.refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  JobBored daily refresh — macOS launchd template.

  Placeholders replaced at install time by scripts/install-launchd-refresh.mjs:
    {{NODE_PATH}}     Node executable path
    {{SCRIPT_PATH}}   scripts/run-scheduled-discovery.mjs
    {{PORT}}          Local worker port
    {{SHEET_ID}}      Optional sheet id override
    {{LOG_PATH}}      Absolute path to the launchd stderr/stdout log
    {{HOUR}}          Hour of day to fire (0-23, local time)
    {{MINUTE}}        Minute of hour to fire (0-59)

  Behavior: builds a fresh command-center.discovery webhook payload at fire
  time and POSTs it to the local browser-use-discovery worker's /webhook
  endpoint. The worker must be running for this to land.

  RunAtLoad=false — loading the plist does NOT immediately fire a refresh.
  To fire manually for testing: `launchctl start com.jobbored.refresh`.
-->
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.jobbored.refresh</string>

    <key>ProgramArguments</key>
    <array>
      <string>{{NODE_PATH}}</string>
      <string>{{SCRIPT_PATH}}</string>
      <string>--trigger</string>
      <string>scheduled-local</string>
      <string>--port</string>
      <string>{{PORT}}</string>
      <string>--sheet-id</string>
      <string>{{SHEET_ID}}</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>{{HOUR}}</integer>
      <key>Minute</key>
      <integer>{{MINUTE}}</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>{{LOG_PATH}}</string>

    <key>StandardErrorPath</key>
    <string>{{LOG_PATH}}</string>
  </dict>
</plist>
